# deploy/growhf.conf
[program:growhf-bot]
command=/opt/growhf/venv/bin/python growhf_reactive_bot.py
directory=/opt/growhf
user=growhf
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/growhf/bot.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
startsecs=10
stopasgroup=true
killasgroup=true
environment=PATH="/opt/growhf/venv/bin"
